# tools/yeydisp-watch.service
[Unit]
Description=yeydisp lab reachability + heap watcher
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=compulab
ExecStart=/home/compulab/yeydisp-signalk/tools/yeydisp-watch.sh
Restart=always
RestartSec=15
# Don't spam the journal with every probe line — they already live in
# /var/log/yeydisp/watch.log. Keep stderr in case the script itself
# misbehaves.
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
